# minios-live-master/linux-live/scripts/10-codexsystem/rootcopy-install/etc/xdg/autostart/codexsystem.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=CodexSystem AI
GenericName=AI Desktop Environment
Comment=Autostart CodexOS Desktop Environment
Exec=/usr/local/bin/codexsystem-launch
Icon=codexsystem
Terminal=false
Categories=Network;AI;Office;
X-GNOME-Autostart-enabled=true
